# configs/nyrelika.service
[Unit]
Description=Nyrelika
[Service]
ExecStart=/usr/bin/python3 /opt/nyrelika/phantom_init.pyw
Restart=always
[Install]
WantedBy=multi-user.target
